# ops/systemd/weitesting-postgres-backup.timer
[Unit]
Description=Run WeiTesting PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
Unit=weitesting-postgres-backup.service

[Install]
WantedBy=timers.target
